# crontab
# Open Brain scheduled jobs (supercronic format)
# Times in UTC

# Learning library selection - 04:30 UTC (before pulse; injects 2-3 learning todos)
30 4 * * * python -m src.jobs.learning_daily

# Morning pulse - 05:00 UTC (08:00 Europe/Vilnius summer)
0 5 * * * python -m src.jobs.pulse

# Daily importance aggregation - 01:00 UTC
0 1 * * * python -m src.jobs.importance

# Commitment miss detection - 00:30 UTC (03:30 Europe/Vilnius summer)
30 0 * * * python -m src.jobs.commitment_miss

# Weekly synthesis - Sunday 00:00 UTC
0 0 * * 0 python -m src.jobs.synthesis

# Weekly training summary - Monday 01:00 UTC (after synthesis)
0 1 * * 1 python -m src.jobs.training_weekly
